# extra/picam_gc/picam_gc.timer
[Unit]
Description=Delete older images

[Timer]
OnBootSec=15min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
